# profile.d/nodejs.sh
# default NODE_ENV to production
export NODE_ENV=${NODE_ENV:production}
